# setup/cockpit.service
[Unit]
Description=JARVIS PI3 — unified phone control server
After=multi-user.target network.target

[Service]
Type=simple
User=root
ExecStartPre=/bin/sleep 5
ExecStart=/bin/bash /home/peepo/PI3/projects/cockpit/start.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
