# setup/cockpit.service
[Unit]
Description=Pi Cockpit — unified phone control server
After=graphical.target media-peepo-KINGSTON.mount
Requires=media-peepo-KINGSTON.mount

[Service]
Type=simple
ExecStartPre=/bin/sleep 8
ExecStart=/bin/bash /media/peepo/KINGSTON/PI3/projects/cockpit/start.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
